# docker/apache/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apachectl -k start -D FOREGROUND
redirect_stderr=true
