# docker/transmission/supervisord.conf
[supervisord]
nodaemon=true

[program:transmission-daemon]
command=/usr/bin/transmission-daemon --logfile=/var/log/transmission.log --log-debug --foreground --config-dir=/etc/transmission
