# apache/letsencrypt/cron.daily/letsencrypt
#!/bin/bash
/usr/bin/cronic /home/mkelly/git/experiments/apache/letsencrypt/letsencrypt_home.michaelkelly.org
